# roles/keepalived/templates/keepalived-exporter.service.j2
[Unit]

[Service]
Restart=always
User=root
ExecStart=/usr/local/bin/keepalived-exporter $ARGS

[Install]
WantedBy=multi-user.target